# roles/smallstep_cli/templates/smallstep-cli-cleanup.timer.j2
[Unit]
Description=Timer for Smallstep Certificate Cleanup

[Timer]
OnCalendar=hourly
RandomizedDelaySec=5m
Persistent=true

[Install]
WantedBy=timers.target
